# Installer/data/SLForceFPS/com.osy86.SLForceFPS.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.osy86.SLForceFPS</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/SLForceFPS</string>
		<string>-r</string>
		<string>20</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>
